# src/router/asterisk/contrib/init.d/rc.gentoo.asterisk
#!/sbin/runscript
# $Id: rc.gentoo.asterisk 67061 2007-06-04 17:11:43Z tilghman $

depend() {
	need net logger
}

start() {
	ebegin "Starting Asterisk"
	/usr/sbin/asterisk
	eend $? "Failed to start Asterisk"
}

stop() {
	ebegin "Stopping Asterisk"
	kill $(cat /var/run/asterisk.pid)
	eend $? "Failed to stop Asterisk"
}
