# src/router/dnsmasq/debian/postrm
#!/bin/sh
set -e

if [ purge = "$1" ]; then
   update-rc.d dnsmasq remove >/dev/null
   if [ -x "$(command -v deluser)" ]; then
     deluser --quiet --system dnsmasq > /dev/null || true
  else
     echo >&2 "not removing dnsmasq system account because deluser command was not found"
  fi
fi
