# src/router/olsrd/openbsd/olsrd.rcd
#!/bin/sh

daemon="/usr/local/sbin/olsrd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
