# src/router/nfs-utils/systemd/nfs-mountd.service
[Unit]
Description=NFS Mount Daemon
DefaultDependencies=no
Requires=proc-fs-nfsd.mount
Wants=network-online.target
After=proc-fs-nfsd.mount
After=network-online.target local-fs.target
After=rpcbind.socket
BindsTo=nfs-server.service

[Service]
Type=forking
ExecStart=/usr/sbin/rpc.mountd
