# src/router/nfs-utils/systemd/rpc_pipefs.target.in
[Unit]
Requires=@_rpc_pipefsmount@
After=@_rpc_pipefsmount@
