# src/router/nfs-utils/systemd/rpc_pipefs.target
[Unit]
Requires=var-lib-nfs-rpc_pipefs.mount
After=var-lib-nfs-rpc_pipefs.mount
