# src/router/nfs-utils/systemd/var-lib-nfs-rpc_pipefs.mount.in
[Unit]
Description=RPC Pipe File System
DefaultDependencies=no
After=systemd-tmpfiles-setup.service
Conflicts=umount.target

[Mount]
What=sunrpc
Where=@_statedir@/rpc_pipefs
Type=rpc_pipefs
