# src/router/opennds/resources/opennds.service
[Unit]
Description=openNDS Captive Portal
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/opennds $OPTIONS
Restart=on-failure

[Install]
WantedBy=multi-user.target
