# src/router/smbd-next/tools/scripts/smbd.service
[Unit]
Description=SMBD userspace service

[Service]
Type=forking
ExecStart=/usr/local/sbin/smbd -s

[Install]
WantedBy=multi-user.target
