# src/router/aiccu/netbsd/PLIST
@comment $NetBSD: PLIST,v 1.1.1.1 2005/12/29 14:53:52 ghen Exp $
man/man1/aiccu.1
sbin/aiccu
share/doc/aiccu/HOWTO
share/doc/aiccu/LICENSE
share/doc/aiccu/README
share/examples/aiccu/aiccu.conf
share/examples/rc.d/aiccu
@dirrm share/examples/aiccu
@dirrm share/doc/aiccu

