# src/router/mstp/utils/mstpd.service.in
[Unit]
Description=MSTP daemon
Before=network-pre.target

[Service]
Type=forking
ExecStart=@bridgestpfile@ restart
PIDFile=@mstpdpidfile@
ExecReload=@bridgestpfile@ restart_config
Restart=always
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
