# src/router/xfsprogs/scrub/xfs_scrub_all.cron.in
10 3 * * 0 root test -e /run/systemd/system || @sbindir@/xfs_scrub_all
