# src/router/xfsprogs/scrub/xfs_scrub_fail@.service.in
[Unit]
Description=Online XFS Metadata Check Failure Reporting for %I
Documentation=man:xfs_scrub(8)

[Service]
Type=oneshot
Environment=EMAIL_ADDR=root
ExecStart=@pkg_lib_dir@/@pkg_name@/xfs_scrub_fail "${EMAIL_ADDR}" %I
User=mail
Group=mail
SupplementaryGroups=systemd-journal
