# src/router/smbd/tools-glib.61/ksmbd.service.in
[Unit]
Description=ksmbd userspace daemon
Requires=modprobe@ksmbd.service
Wants=network-online.target
After=modprobe@ksmbd.service network.target network-online.target

[Service]
Type=forking
PIDFile=@runstatedir@/ksmbd.lock
ExecStart=@sbindir@/ksmbd.mountd
ExecReload=@sbindir@/ksmbd.control --reload
ExecStop=@sbindir@/ksmbd.control --shutdown

[Install]
WantedBy=multi-user.target
