# src/router/gpsd/systemd/gpsd.service
[Unit]
Description=GPS (Global Positioning System) Daemon
Requires=gpsd.socket

[Service]
ExecStart=/usr/sbin/gpsd -N

[Install]
Also=gpsd.socket
