# src/router/lvm2/scripts/dm_event_systemd_red_hat.socket.in
[Unit]
Description=Device-mapper event daemon FIFOs
Documentation=man:dmeventd(8)
DefaultDependencies=no

[Socket]
ListenFIFO=@DEFAULT_DM_RUN_DIR@/dmeventd-server
ListenFIFO=@DEFAULT_DM_RUN_DIR@/dmeventd-client
SocketMode=0600
RemoveOnStop=true

[Install]
WantedBy=sockets.target
