# src/router/mdadm/udev-md-clustered-confirm-device.rules
# do not edit this file, it will be overwritten on update

SUBSYSTEM!="block", GOTO="clustermd_end"

# handle md arrays
KERNEL!="md*", GOTO="clustermd_end"
ENV{DEVTYPE}!="disk", GOTO="clustermd_end"
ACTION!="change", GOTO="clustermd_end"
ENV{EVENT}!="ADD_DEVICE", GOTO="clustermd_end"
ENV{DEVICE_UUID}!="?*", GOTO="clustermd_end"
ENV{RAID_DISK}!="?*", GOTO="clustermd_end"

# Based on the received UUID, node confirms the device if
# it is found by blkid, otherwise the node reports it is
# missing.
PROGRAM="BINDIR/blkid -o device -t UUID_SUB=$env{DEVICE_UUID}", ENV{.md.newdevice} = "$result"

ENV{.md.newdevice}!="", RUN+="BINDIR/mdadm --manage $env{DEVNAME} --cluster-confirm $env{RAID_DISK}:$env{.md.newdevice}"
ENV{.md.newdevice}=="", RUN+="BINDIR/mdadm --manage $env{DEVNAME} --cluster-confirm $env{RAID_DISK}:missing"

LABEL="clustermd_end"
